# pyfreqd.service
[Unit]
Description=pyfreqd daemon

[Service]
ExecStart=/usr/local/sbin/pyfreqd

[Install]
WantedBy=multi-user.target
